# server/local.d/01_tailscale_nokeys.start
#!/bin/sh
tailscale up \
    --hostname alpine \
    --auth-key  tskey-auth-secret-key

# nginx is dependent on tailscale hosted machine
# it'd [emerg] error on start
#
# once tailscale network is ready, start nginx again
rc-service nginx start
