# server/local.d/02_mjpeg-streamer.start
#!/bin/sh
mjpg_streamer -i "input_uvc.so --resolution 176x144"
