# server/local.d/03_mjpeg-streamer.start
#!/bin/sh
ffmpeg -f v4l2 -framerate 30 -video_size 176x144 -input_format yuyv422 -i /dev/video0 -ar 44100 -ac
