# server/local.d/tailscale.start
#!/bin/sh
tailscale up \
	--hostname alpine \
	--auth-key  tskey-auth-kwJFr5HoF911CNTRL-fjjxnYNzr7MLex4bivbG8McvnvrT7MBnQ
