# server/local.d/tailscale.start
#!/bin/sh
tailscale up \
    --hostname CHANGE_ME \
	--auth-key  tskey-auth-CHANGE-ME-WITH-REUSABLE-AUTH-KEY
