# server/local.d/tailscale.start
#!/bin/sh
tailscale up \
	--auth-key  tskey-auth-CHANGE-ME-WITH-REUSABLE-AUTH-KEY
